# docker/s6-rc.d/qcluster/run
#!/command/with-contenv sh
cd /var/www/seismic/api/

python3 manage.py qcluster
