# docker/s6-rc.d/qcluster/run
#!/command/with-contenv sh
cd /var/www/seismic/backend/ || exit 1

python3 manage.py qcluster
